# crontab
10 0 * * * /home/point/point/stat/run.sh
